# config/weston/my/weston-autostart.service
[Unit]
Description=Start a wayland application

[Service]
Type=simple
ExecStart=/opt/sbin/autostartd
