# config/autostart/my/autostart.service
[Unit]
Description=Auto start wm complementary services

[Service]
Type=simple
ExecStart=/opt/sbin/autostartd
