# config/autostart/audio.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "sleep 2; ~/scripts/sys_mon/control.lua fun pa_set_default;"
Name=Audio
Comment=Autostart audio at login
