# config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "sleep 2; dunst;"
Name=dunst
Comment=Autostart dunst at login
