# config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "sleep 1; picom;"
Name=picom
Comment=Autostart picom at login
